# odoo_v9/etc/supervisord/odoo.conf
[program:odoo]
user=root
group=root
command=/usr/local/bin/run_odoo.sh
autostart=true
autorestart=true
